# bin/cron/cron.hourly/rsync-elibrary
#!/bin/bash
! mountpoint -q /media/binz/backup/ || rsync -avz --progress ./ /media/binz/backup/elibrary/ --delete
