# bin/cron/cron.hourly/rsync-dev
#!/bin/bash
! mountpoint -q /media/binz/backup/ || rsync -aAXv ~/dev --exclude=node_modules /media/binz/backup/rsync/dev --delete

